# deploy/configs/systemd/rainbow-say-backend.service
[Unit]
Description=Rainbow Say Django Backend
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/rainbow-say/backend
Environment="PATH=/var/www/rainbow-say/backend/venv/bin"
EnvironmentFile=/var/www/rainbow-say/backend/.env
ExecStart=/var/www/rainbow-say/backend/venv/bin/gunicorn \
    --bind 127.0.0.1:8001 \
    --workers 3 \
    --timeout 120 \
    --access-logfile /var/log/rainbow-say-backend-access.log \
    --error-logfile /var/log/rainbow-say-backend-error.log \
    config.wsgi:application
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

